# test/test.service
[Unit]
Description=Monitor test process

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test.sh
